# moby-engine/postinstall/deb/prerm

if [ -d /run/systemd/system ] ; then
	systemctl --system daemon-reload >/dev/null || true
fi

if [ -x "/usr/bin/deb-systemd-helper" ]; then
    deb-systemd-helper mask 'docker.service' 'docker.socket' >/dev/null || true
fi

if [ -z "${DPKG_ROOT:-}" ] && [ -d /run/systemd/system ] ; then
	deb-systemd-invoke stop 'docker.service' 'docker.socket' >/dev/null || true
fi
                    