# moby-containerd-shim-systemd/postinstall/deb/prerm

if [ -d /run/systemd/system ] ; then
	systemctl --system daemon-reload >/dev/null || true
fi

if [ -x "/usr/bin/deb-systemd-helper" ]; then
    deb-systemd-helper mask 'containerd-shim-systemd-v1.service' 'containerd-shim-systemd-v1.socket' >/dev/null || true
fi

if [ -z "${DPKG_ROOT:-}" ] && [ -d /run/systemd/system ] ; then
	deb-systemd-invoke stop 'containerd-shim-systemd-v1.service' 'containerd-shim-systemd-v1.socket' >/dev/null || true
fi
                    